# supervisor/sshd.conf
[program:sshd]
directory=/
command=/usr/sbin/sshd -f /app/code/sshd_config -D
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

